# supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/app/logs

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# ── PostgreSQL ────────────────────────────────────────────────────────────────
[program:postgres]
command=/bin/bash -c "PGV=$(ls /usr/lib/postgresql/ | sort -V | tail -1) && gosu postgres /usr/lib/postgresql/$PGV/bin/postgres -D /var/lib/postgresql/data"
autostart=true
autorestart=true
priority=10
stdout_logfile=/app/logs/postgres.log
stderr_logfile=/app/logs/postgres.log

# ── Redis ─────────────────────────────────────────────────────────────────────
[program:redis]
command=/bin/bash -c "redis-server --requirepass \"$(grep REDIS_PASSWORD /app/.env | cut -d= -f2 | tr -d ' ')\" --bind 127.0.0.1"
autostart=true
autorestart=true
priority=10
stdout_logfile=/app/logs/redis.log
stderr_logfile=/app/logs/redis.log

[program:xray]
command=/usr/local/bin/xray run -config /etc/xray/config.json
autostart=true
autorestart=true
startretries=10
priority=15
stdout_logfile=/app/logs/xray.log
stderr_logfile=/app/logs/xray.log

[program:xray_node]
command=uvicorn xray_node.main:app --host 0.0.0.0 --port 8090
directory=/app
autostart=true
autorestart=true
startretries=10
priority=14
stdout_logfile=/app/logs/xray_node.log
stderr_logfile=/app/logs/xray_node.log

# ── DB Init (runs once) ───────────────────────────────────────────────────────
[program:db_init]
command=python /app/scripts/db_init.py
directory=/app
autostart=true
autorestart=false
startsecs=0
startretries=3
priority=20
stdout_logfile=/app/logs/db_init.log
stderr_logfile=/app/logs/db_init.log

# ── Telegram Bot ──────────────────────────────────────────────────────────────
[program:bot]
command=python -m bot.main
directory=/app
autostart=false
autorestart=true
priority=30
stdout_logfile=/app/logs/bot.log
stderr_logfile=/app/logs/bot.log

# ── Subscription Server ───────────────────────────────────────────────────────
[program:subscription]
command=uvicorn subscription_server.main:app --host 127.0.0.1 --port 8000
directory=/app
autostart=true
autorestart=true
priority=30
stdout_logfile=/app/logs/subscription.log
stderr_logfile=/app/logs/subscription.log

# ── Website API ───────────────────────────────────────────────────────────────
[program:website]
command=uvicorn website.api:app --host 127.0.0.1 --port 8002
directory=/app
autostart=true
autorestart=true
priority=30
stdout_logfile=/app/logs/website.log
stderr_logfile=/app/logs/website.log

# ── Admin Panel ───────────────────────────────────────────────────────────────
[program:admin_panel]
command=uvicorn admin_panel.main:app --host 127.0.0.1 --port 8001
directory=/app
autostart=true
autorestart=true
priority=30
stdout_logfile=/app/logs/admin_panel.log
stderr_logfile=/app/logs/admin_panel.log

# ── Celery ────────────────────────────────────────────────────────────────────
[program:celery]
command=celery -A utils.celery_app worker --beat -l info
directory=/app
autostart=true
autorestart=true
priority=30
stdout_logfile=/app/logs/celery.log
stderr_logfile=/app/logs/celery.log

# ── Nginx ─────────────────────────────────────────────────────────────────────
[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
priority=40
stdout_logfile=/app/logs/nginx.log
stderr_logfile=/app/logs/nginx.log

# ── Certbot Auto-Renew (каждые 12 часов) ─────────────────────────────────────
[program:certbot_renew]
command=/bin/bash -c "while true; do sleep 43200; certbot renew --nginx --quiet 2>/app/logs/certbot_renew.log; done"
autostart=true
autorestart=true
priority=50
stdout_logfile=/app/logs/certbot.log
stderr_logfile=/app/logs/certbot.log

[program:support_bot]
command=python -m support_bot.main
directory=/app
autostart=true
autorestart=true
priority=30
stdout_logfile=/app/logs/support_bot.log
stderr_logfile=/app/logs/support_bot.log
